# tools/sshd-localhost
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /home/cvs/kierbiischt-ion-sysadmin/tools/sshd-localhost,v 1.1 2005/05/09 12:05:41 steve Exp $

depend() {
	use logger dns
	need net
}

checkconfig() {
	if [ ! -d /var/empty ] ; then
		mkdir -p /var/empty || return 1
	fi

	if [ ! -e /etc/ssh/sshd_config.localhost ] ; then
		eerror "You need an /etc/ssh/sshd_config.localhost file to run sshd localhost shim"
		return 1
	fi

	gen_keys
}

gen_keys() {
	if [ ! -e /etc/ssh/ssh_host_key ] ; then
		einfo "Generating Hostkey..."
		/usr/bin/ssh-keygen -t rsa1 -b 1024 -f /etc/ssh/ssh_host_key -N ''
	fi
	if [ ! -e /etc/ssh/ssh_host_dsa_key ] ; then
		einfo "Generating DSA-Hostkey..."
		/usr/bin/ssh-keygen -d -f /etc/ssh/ssh_host_dsa_key -N ''
	fi
	if [ ! -e /etc/ssh/ssh_host_rsa_key ] ; then
		einfo "Generating RSA-Hostkey..."
		/usr/bin/ssh-keygen -t rsa -f /etc/ssh/ssh_host_rsa_key -N ''
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting sshd localhost shim"
	start-stop-daemon --start --quiet --pidfile /var/run/sshd-localhost.pid \
		--startas /usr/sbin/sshd -- ${SSHD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping sshd localhost shim"
	start-stop-daemon --stop --quiet --pidfile /var/run/sshd-localhost.pid
	eend $?
}
